# borg-backup-script-example.timer
[Unit]
Description=Create borg backup for 

[Timer]
OnCalendar=
OnBootSec=2min
OnUnitActiveSec=5min


[Install]
WantedBy=multi-user.target